# infra/docling-service/docling.service
[Unit]
Description=Docling Converter Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/docling-service
ExecStart=/opt/docling-service/.venv/bin/python main.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/docling-service/.env

[Install]
WantedBy=multi-user.target
